# root_files/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/sbin/getty -t 0 -L -n tty_console 115200 xterm
::shutdown:/bin/umount -a -r

